# data/supervisor.conf
[group:fetch_items]
programs=fetch_items_1,fetch_items_2,fetch_items_3

[program:fetch_items_1]
directory=/path/to/f43
command=/usr/bin/php bin/console swarrot:consume:f43.fetch_items --env=prod f43.fetch_items
autostart=true
autorestart=true
stderr_logfile=/space/logs/supervisord/fetch_items_1.err
stdout_logfile=/space/logs/supervisord/fetch_items_1.log
user=deploy
environment = http_proxy="",https_proxy=""

[program:fetch_items_2]
directory=/path/to/f43
command=/usr/bin/php bin/console swarrot:consume:f43.fetch_items --env=prod f43.fetch_items
autostart=true
autorestart=true
stderr_logfile=/space/logs/supervisord/fetch_items_2.err
stdout_logfile=/space/logs/supervisord/fetch_items_2.log
user=deploy
environment = http_proxy="",https_proxy=""

[program:fetch_items_3]
directory=/path/to/f43
command=/usr/bin/php bin/console swarrot:consume:f43.fetch_items --env=prod f43.fetch_items
autostart=true
autorestart=true
stderr_logfile=/space/logs/supervisord/fetch_items_3.err
stdout_logfile=/space/logs/supervisord/fetch_items_3.log
user=deploy
environment = http_proxy="",https_proxy=""
