# data/supervisor.conf
[group:fetch_items]
programs=fetch_items_1,fetch_items_2,fetch_items_3

[program:fetch_items_1]
directory=/path/to/f43
command=/usr/bin/php bin/console messenger:consume --limit=20 --env=prod fetch_items
autostart=true
autorestart=true
stderr_logfile=/path/to/f43/var/log/fetch_items_1.err
stdout_logfile=/path/to/f43/var/log/fetch_items_1.log
user=www-data
environment = http_proxy="",https_proxy=""

[program:fetch_items_2]
directory=/path/to/f43
command=/usr/bin/php bin/console messenger:consume --limit=20 --env=prod fetch_items
autostart=true
autorestart=true
stderr_logfile=/path/to/f43/var/log/fetch_items_2.err
stdout_logfile=/path/to/f43/var/log/fetch_items_2.log
user=www-data
environment = http_proxy="",https_proxy=""

[program:fetch_items_3]
directory=/path/to/f43
command=/usr/bin/php bin/console messenger:consume --limit=20 --env=prod fetch_items
autostart=true
autorestart=true
stderr_logfile=/path/to/f43/var/log/fetch_items_3.err
stdout_logfile=/path/to/f43/var/log/fetch_items_3.log
user=www-data
environment = http_proxy="",https_proxy=""
